# Crontab/psn.crontab
0 16 * * * cd ~/psn-discount/Crontab && ./psn.sh > psn.log 2>&1
